# kit/init/systemd-templates/compai-merch.service.tmpl
[Unit]
Description=Compai @BRAND@ · merch agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@USER@
Group=@USER@
WorkingDirectory=@HOME@/agents/merch
Environment="COMPAI_HOME=@HOME@"
Environment="COMPAI_BRAND=@BRAND@"
Environment="COMPAI_AGENT=merch"
ExecStart=/usr/bin/python3 @HOME@/services/agent-runner.py --agent merch
Restart=on-failure
RestartSec=30
StandardOutput=append:@HOME@/logs/merch.log
StandardError=append:@HOME@/logs/merch.log

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=@HOME@/logs @HOME@/brain @HOME@/agents/merch /tmp

[Install]
WantedBy=multi-user.target
