# kit/init/systemd-templates/operai-mcp.service.tmpl
[Unit]
Description=OperAI @BRAND@ · MCP server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@USER@
Group=@USER@
WorkingDirectory=@HOME@/services/mcp
Environment="OPERAI_HOME=@HOME@"
Environment="OPERAI_BRAND=@BRAND@"
ExecStart=/usr/bin/python3 @HOME@/services/mcp/server.py --host 127.0.0.1 --port 8787
Restart=on-failure
RestartSec=10
StandardOutput=append:@HOME@/logs/mcp.log
StandardError=append:@HOME@/logs/mcp.log

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=@HOME@/logs @HOME@/brain @HOME@/services/mcp

[Install]
WantedBy=multi-user.target
