# kit/init/systemd-templates/operai-factory-runtime.service.tmpl
[Unit]
Description=Compai @BRAND@ · factory runtime daemon
After=network-online.target compai-mcp.service
Wants=network-online.target

[Service]
Type=simple
User=@USER@
Group=@USER@
WorkingDirectory=@HOME@
Environment="COMPAI_HOME=@HOME@"
Environment="COMPAI_BRAND=@BRAND@"
Environment="PYTHONPATH=@HOME@/services/init/cli:@HOME@/workflows"
ExecStart=/usr/bin/python3 -m compai_init.factory_runtime.daemon --home @HOME@ --brand @BRAND@
Restart=on-failure
RestartSec=10
StandardOutput=append:@HOME@/logs/factory-runtime.log
StandardError=append:@HOME@/logs/factory-runtime.log

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=@HOME@/events @HOME@/brain @HOME@/logs @HOME@/agents /tmp

[Install]
WantedBy=multi-user.target
